# keyboard_touchpad_helper/keyboard_touchpad_helper.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Keyboard Monitor (helps touchpad palm detection)"
author          "chromium-os-dev@chromium.org"

# This daemon maintains a connection to the X server, thus it uses non-standard
# stop/stop settings:
start on started ui
stop on stopping ui

respawn
expect fork

env XAUTHORITY=/home/chronos/.Xauthority
env DISPLAY=:0.0

exec /usr/sbin/keyboard_touchpad_helper \
  --tp_prop_name_high "Keyboard Touched Timeval High" \
  --tp_prop_name_low "Keyboard Touched Timeval Low"
